# powerfail-proxmox.service
[Unit]
Description=UPS Power Failure Check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/powerfail-proxmox.sh
